# com.user.githealthcheck.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.githealthcheck</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/anthonygarand/Projects/ai-assistant/git-health-check.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/anthonygarand/Projects/ai-assistant/logs/git-health-check.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/anthonygarand/Projects/ai-assistant/logs/git-health-check-error.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
